# cron/crontab
*/5 * * * * root /usr/bin/php -f /var/www/html/cron.php
